# sv/ratpoison/run
#!/bin/sh
sleep 5
exec chpst -u $USER ratpoison -d :0 -f ./ratpoison.conf
