# sv/x/run
#!/bin/sh
exec X :0 -s 0 -dpms -nolock -nocursor
